# com.user.newsserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.newsserver</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kartiknath/Documents/code/code/ai_projects/news-aggregator/venv/bin/python</string>
        <string>-m</string>
        <string>http.server</string>
        <string>8080</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/kartiknath/Documents/code/code/ai_projects/news-aggregator/output</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/kartiknath/Documents/code/code/ai_projects/news-aggregator/logs/server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kartiknath/Documents/code/code/ai_projects/news-aggregator/logs/server-error.log</string>
</dict>
</plist>
